# supervisor.deluge.ini
[program:deluged]
command=/usr/bin/deluged --config /mnt/deluge/config -d
autostart=true
autorestart=true

[program:deluge-web]
command=/usr/bin/deluge-web --config /mnt/deluge/config
autostart=true
autorestart=true
